# services/autorestic/systemd/autorestic-cron.service
[Unit]
Description=autorestic cron service

[Service]
Type=oneshot
ExecStart=autorestic --ci cron
WorkingDirectory=%h
